# cookbooks/ey-puma/templates/default/puma.service.erb
[Unit]
Description=<%= @app %> Puma
After=network.target

Requires=puma_<%= @app %>.socket

[Service]
Type=notify

EnvironmentFile=/data/<%= @app %>/shared/config/env

<% if @systemctlvar %>
EnvironmentFile=/data/<%= @app %>/shared/config/env.systemctl
<% end %>
<% if @customvar %>
EnvironmentFile=/data/<%= @app %>/shared/config/env.custom
<% end %>

WatchdogSec=60

MemoryMax=<%= @app_memory_limit %>M


User=<%= @username %>

WorkingDirectory=/data/<%= @app %>/current

ExecStartPre=/usr/bin/sudo /usr/bin/chown -R <%= @username %>:<%= @username %> /run/engineyard/<%= @app %>

ExecStart=/data/<%= @app %>/current/ey_bundler_binstubs/puma -w <%= @workers %>:<%= @threads %> -e <%= @framework_env %> --bind unix:///var/run/engineyard/<%= @app %>/puma_<%= @app %>.sock --port <%= @port %> --control-url unix:///var/run/engineyard/<%= @app %>/puma_<%= @app %>-ctl.sock --state /var/run/engineyard/<%= @app %>/puma_<%= @app %>.state --dir /data/<%= @app %>/current/ --prune-bundler

Restart=always

[Install]
WantedBy=multi-user.target
